# net-proxy/frp/files/frpc.service
[Unit]
Description=frpc daemon

[Service]
Type=simple
ExecStart=/usr/bin/frpc -c /etc/frp/frpc.ini

[Install]
WantedBy=multi-user.target
